# packages/adsms/adsms-0.5.tar.gz/adsms-0.5/debian/adsms.service
# Systemd unit file for the adsms service
# Generated by ChatGPT

[Unit]
Description=Send SMS notifications based on ADS-B data
After=network.target
Requires=tar1090.service
After=tar1090.service

[Service]
User=adsms
Group=adsms
ExecStart=/usr/local/bin/adsms /etc/adsms.json
Restart=always

[Install]
WantedBy=multi-user.target
